# docs/systemd/flight-server.service
[Unit]
Description=BotCopier Arrow Flight Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/BotCopier
ExecStart=/usr/bin/env botcopier-flight-server
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
