# docs/systemd/online-trainer.service
[Unit]
Description=BotCopier Online Trainer
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/BotCopier
ExecStart=/usr/bin/python3 scripts/online_trainer.py --flight 127.0.0.1:8815
Restart=on-failure
WatchdogSec=60
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
