# docs/systemd/promote-best-models.timer
[Unit]
Description=Schedule Promotion of Best Models

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
