# docs/systemd/promote-best-models.service
[Unit]
Description=BotCopier Promote Best Models
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/BotCopier
ExecStart=/usr/bin/python3 scripts/promote_best_models.py models best --metric sharpe
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
